# traffic_platform/systemd/traffic-platform.service
[Unit]
Description=Traffic Event Web Platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=firefly
WorkingDirectory=/home/firefly/dds_test/camera_radar_fusion_process
Environment=TRAFFIC_PLATFORM_PORT=8080
ExecStart=/usr/bin/python3 /home/firefly/dds_test/camera_radar_fusion_process/traffic_platform/backend/app.py --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
